# scripts/termtile.service
[Unit]
Description=Terminal Tiling Manager
Documentation=https://github.com/1broseidon/termtile
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/termtile daemon
Restart=on-failure
RestartSec=5

# Environment (optional; generated by install script)
EnvironmentFile=-%h/.config/termtile/env

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
